# assets/move-appimages.service
[Unit]
Description=Move user-downloaded AppImages to central dir


[Service]
Type=oneshot
ExecStart=/usr/local/bin/appiman ingest