# assets/move-appimages.timer
[Unit]
Description=Run AppImage ingestion every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=1s

[Install]
WantedBy=timers.target
