# assets/move-appimages.path
[Unit]
Description=Watch user home dirs for AppImage downloads

[Path]
PathChanged=/home
PathModified=/home

[Install]
WantedBy=multi-user.target
