# assets/register-appimages.path
[Unit]
Description=Watch for new AppImages in /opt/applications/raw


[Path]
PathChanged=/opt/applications/raw
PathModified=/opt/applications/raw
DirectoryNotEmpty=/opt/applications/raw


[Install]
WantedBy=multi-user.target