# assets/register-appimages.service
[Unit]
Description=Register AppImages system-wide


[Service]
Type=oneshot
ExecStart=/usr/local/sbin/register-appimages.sh