# assets/register-appimages.path
[Unit]
Description=Watch for new AppImages in /opt/applications/raw


[Path]
PathExistsGlob=/opt/applications/raw/*.AppImage
PathExistsGlob=/opt/applications/raw/*.appimage


[Install]
WantedBy=multi-user.target
